# apt/gpm.service
[Unit]
Description=Go Performance Monitor

[Service]
ExecStart=/etc/gpm/main
Restart=always
User=root

[Install]
WantedBy=multi-user.target